# ansible/roles/node_app/templates/node-app.service.j2
[Unit]
Description=My Node.js App
After=network.target [cite: 112, 113, 114]

[Service]
User=deploy
WorkingDirectory=/home/deploy/apps/myapp
ExecStart=/usr/bin/node /home/deploy/apps/myapp/index.js
Restart=always
Environment=NODE_ENV=production [cite: 115, 116, 117, 118, 119, 120, 124]

[Install]
WantedBy=multi-user.target [cite: 122, 123]